# projects/ollama/scripts/com.local.usage-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.local.usage-monitor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>SCRIPT_PATH_PLACEHOLDER</string>
	</array>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>RunAtLoad</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>USAGE_MONITOR_MEM_PCT</key>
		<string>90</string>
		<key>USAGE_MONITOR_CPU_PCT</key>
		<string>95</string>
	</dict>
</dict>
</plist>
